# dyndns.timer
[Unit]
Description=timer for dyndns.service

[Timer]
OnBootSec=10s
OnUnitActiveSec=15m
Unit=dyndns.service

[Install]
WantedBy=timers.target
